# 4-a/etc/supervisor/conf.d/backup.conf
[program:backup]
command=whenever --update-crontab -f /root/Backup/schedule.rb
startretries=0
